# data/eos-app-manager-migrate-bundles.service.in
[Unit]
Description=Migrate old application bundle installations
Documentation=man:eamd(8)
DefaultDependencies=false
Conflicts=shutdown.target
After=local-fs.target systemd-journald.socket
Before=systemd-update-done.service shutdown.target graphical.target
ConditionNeedsUpdate=/var
RequiresMountsFor=/var/endless /var/endless-extra

[Service]
Type=oneshot
ExecStart=@bindir@/eamctl ensure-symlink-farm
User=@EAM_USER@

[Install]
WantedBy=multi-user.target
