# data/eos-app-manager.service.in
[Unit]
Description=Endless App Manager
Documentation=man:eamd(8)

[Service]
Type=dbus
BusName=com.endlessm.AppManager
ExecStart=@bindir@/eamd
User=@EAM_USER@
